# zfs/files/zpool-scrub@.service
[Unit]
Description = zpool scrub on %i

[Service]
ExecStart = zpool scrub -w %i
