# zfs/files/zpool-trim@.service
[Unit]
Description = zpool trim on %i

[Service]
ExecStart = zpool trim -w %i
